# deploy/furnace.service
[Unit]
Description=Furnace Deployment Agent
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=furnace
Group=furnace
ExecStart=/usr/local/bin/furnace --config /etc/furnace/furnace.yaml
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/lib/furnace /srv/apps
ReadOnlyPaths=/etc/furnace

# Docker socket access via group membership
SupplementaryGroups=docker

# Resource limits
LimitNOFILE=65536
MemoryMax=256M

[Install]
WantedBy=multi-user.target
